# keystored.service
[Unit]
Description=Key-Value Storage Daemon
After=network.target

[Service]
Type=simple
User=keystored
Group=keystored
ExecStart=/usr/local/bin/keystored 
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
SyslogIdentifier=keystored
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
